# plists/com.macpilot.health-check.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.macpilot.health-check</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>MACPILOT_AGENT_NAME</key>
    <string>health-check</string>
  </dict>
  <key>ProgramArguments</key>
  <array>
    <string>__MACPILOT_DIR__/agents/health-check.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>9</integer>
    <key>Minute</key><integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>__MACPILOT_DIR__/logs/health-check.out</string>
  <key>StandardErrorPath</key>
  <string>__MACPILOT_DIR__/logs/health-check.err</string>
</dict>
</plist>
